# roles/galera/templates/mariadb.service.j2
[Service]
RuntimeDirectory={{ play_galera.socket | dirname }}